# roles/rpi-mh-z19c-exporter/templates/rpi-mh-z19c-exporter.service.j2
[Unit]
Description=MH-Z19C Prometheus exporter

[Service]
Type=notify
ExecStart={{ binary_path }}
Environment="CO2_SENSOR_PATH={{ sensor_path }}"
Environment="CO2_METRICS_LISTEN_ADDRS={{ listen_addrs | join(' ') }}"
Restart=on-failure
User=rpi-mh-z19c-exporter
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
NoNewPrivileges=yes

[Install]
WantedBy=default.target